# ops/systemd/ttangbu-backend.service
[Unit]
Description=ttangbu backend service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ttangbu/backend
EnvironmentFile=/opt/ttangbu/backend/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
